# systemd/pxBackupManager.service
[Unit]
Description=pxBackupManager - Backup Service for FiveM and MariaDB
After=network.target mysql.service
Wants=pxBackupManager.timer

[Service]
Type=oneshot
User=pxbackup
Group=pxbackup
WorkingDirectory=/opt/pxBackupManager
ExecStart=/opt/pxBackupManager/pxBackupManager \
  -fivem /opt/fivem \
  -db-name fivem \
  -db-user backup \
  -db-pass YourPasswordHere \
  -backup-dir /backups/pxBackupManager \
  -retention 30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pxBackupManager

# Security
PrivateTmp=yes
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
