# systemd/pxBackupManager.timer
[Unit]
Description=pxBackupManager - Daily Backup Timer
Requires=pxBackupManager.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=*-*-* 02:00:00

# Run on boot if missed
Persistent=true

# Random delay to prevent system overload (5 minutes max)
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
